# Website/Init/PSE_Website_Server.service
[Unit]
Description=Website Server

[Service]
ExecStart=python /usr/local/bin/website_controller.py

[Install]
WantedBy=multi-user.target