# RaspberryPi/Init/PSE_Website.service
[Unit]
Description=Website Raspberry-Pi

[Service]
ExecStart=python /usr/local/bin/website_pi.py
User=pi

[Install]
WantedBy=multi-user.target